# debian/DEBIAN/prerm
#!/bin/bash

systemctl disable hobot-modules-load.service

exit 0

